# project2/docker/server_img/etc/service/tor/run
#!/bin/sh

# Ensures Tor's HOME has the right owner.
chown -R debian-tor:debian-tor /var/lib/tor

exec 2>&1
exec chpst -P -u debian-tor:debian-tor env HOME=/var/lib/tor tor
